# scripts/systemd/db-backup.timer
[Unit]
Description=Run Database Backup Utility daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

